# rootfs/etc/services.d/miner/run
#!/bin/sh

echo "* generating block..."

litecoin-cli -conf=/etc/conf/node-1/litecoin.conf generate 1
sleep 10